# linux-course/11-scheduling-logging-backups/code/report.timer
[Unit]
Description=Run report.service every 2 minutes (Linux course demo)

[Timer]
OnBootSec=1min          # first run a minute after boot
OnUnitActiveSec=2min    # then every 2 minutes
Persistent=true         # catch up if the machine was off at the scheduled time

[Install]
WantedBy=timers.target
